# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:nuxt]
command=node .output/server/index.mjs
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=1

[program:flask]
command=gunicorn -w 1 -k gthread --threads 2 --timeout 180 --graceful-timeout 45 --max-requests 200 --max-requests-jitter 30 -b 0.0.0.0:5000 app:app
directory=/app/flask
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=2
environment=PYTHONUNBUFFERED=1
